# deploy/jobl-api.service
[Unit]
Description=Jobl API Service
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/home/webadmin/Jobl/api.jobl.ai
EnvironmentFile=/home/webadmin/Jobl/api.jobl.ai/.env
ExecStart=/home/webadmin/Jobl/api.jobl.ai/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
